# deploy/mcp-binance.service
[Unit]
Description=MCP Binance (TypeScript)
After=network.target

[Service]
Type=simple
WorkingDirectory=/opt/mcp-binance
ExecStart=/usr/bin/node dist/index.js
Restart=always
Environment=NODE_ENV=production
Environment=PORT=8080
Environment=HOST=0.0.0.0
Environment=AUTH_TOKEN=change-me
Environment=TRANSPORT=http

[Install]
WantedBy=multi-user.target
